# recipes-fotahub/fotahubclient/files/fotahubd.service.in
[Unit]
Description=FotaHub Client daemon
After=network.target

[Service]
WorkingDirectory=/usr/bin
ExecStart=/usr/bin/python3 /usr/bin/fotahubd --config {{config}} -v
KillMode=process

[Install]
WantedBy=multi-user.target
